# ansible/roles/smbmount/files/foo.automount
[Unit]
Requires=remote-fs-pre.target
After=remote-fs-pre.target

[Install]
WantedBy=remote-fs.target
